# supervisor/conf.d/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
user=www-data
autostart=true
